# debian/bluetracker.postinst
#!/bin/sh

set -e

#DEBHELPER#

update-rc.d bluetracker defaults 98 02
