# debian/bluetracker.postinst
#!/bin/sh

#DEBHELPER#

update-rc.d bluetrackd defaults 98 02 > /dev/null
