# debian/gyrid.prerm
#!/bin/sh

set -e

if [ -f /var/run/bluetracker.pid ]; then
    echo "Stopping bluetracker."

    if which invoke-rc.d >/dev/null 2>&1; then
        invoke-rc.d bluetracker stop
    else
        /etc/init.d/bluetracker stop
    fi
fi

#DEBHELPER#

exit 0
