# debian/gyrid.postinst
#!/bin/sh

set -e

#DEBHELPER#

update-rc.d gyrid defaults 95 02
